# docker/neurofeed.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
CRON_TZ=America/Sao_Paulo

0 7 * * * root set -a && . /run/neurofeed-cron-env.sh && set +a && /usr/local/bin/neurofeed >> /proc/1/fd/1 2>&1
